# stranded.conf
description "Start the stranded.net node.js server"
author "Jonny Balls"

start on runlevel [2345]

exec nodejs /home/stranded/the-stranded/app.js
